# uwsgi/templates/uwsgi_upstart.conf.j2
description "{{uwsgi_service_name}}"

start on runlevel [2345]
stop on runlevel [06]

respawn

exec {{uwsgi_executable}} --die-on-term {{uwsgi_config_directory}}/{{uwsgi_config_file}}
